# svalin/packaging/systemd/svalin-agent.service
[Unit]
Description=Agent for Svalin

[Service]
Type=simple
ExecStart=/usr/bin/svalin agent
Restart=always

[Install]
WantedBy=multi-user.target